# root/etc/s6-overlay/s6-rc.d/init-mylar3-config/up
/etc/s6-overlay/s6-rc.d/init-mylar3-config/run
