# root/etc/services.d/mylar/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python3 /app/mylar3/Mylar.py --nolaunch \
	--datadir /config/mylar
